# docker/docker_utils/cloud_sql.template.conf
[program:cloudsql]
command={{cloud_sql_exec}} -dir {{cloud_sql_dir}} -projects {{google_project}} -credential_file {{credential_file}}
user=root
autostart=true
autorestart=true
