# docker/docker_utils/redis.conf
[program:redis]
command=/usr/local/bin/redis-server /etc/redis.conf
#directory=/storage/redis_data
user=root
autostart=true
autorestart=true
stopsignal=QUIT
